# .config/systemd/user/btrfs-backup.timer
[Unit]
Description=Daily Btrfs Backup After Login

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=default.target
